# rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv sh

exec dnsmasq -d --conf-file=/config/etc/dnsmasq.conf
